# deploy/supervisor/mzcorp-worker.conf
; mzcorp-worker — Laravel queue:work воркеры.
;
; Прод: /etc/supervisor/conf.d/mzcorp-worker.conf
; После правки: sudo supervisorctl reread && sudo supervisorctl update
;                && sudo supervisorctl restart mzcorp-worker:*
;
; Очереди (см. post-mortem 2026-05-28 в MEMORY § «Sync info@ stuck»):
;   mail-sync       — IMAP SyncMailboxFolderJob (time-critical).
;                     Должно крутиться даже если default забит.
;   default         — всё остальное (Parse, Deliver, Route, KB resolve,
;                     OutgoingQuoteParse, нотификации, AI-decisions).
;   catalog-resolve — тяжёлая catalog/LLM/pgvector-работа
;                     (ResolvePendingChunkJob и dispatcher).
;                     Изолирована, чтобы её таймауты не блокировали
;                     mail-sync.
;
; Порядок --queue ВАЖЕН: Laravel queue:work обрабатывает очереди в
; том порядке, в котором они переданы. mail-sync,default,catalog-resolve
; значит «сначала забери всё из mail-sync, потом из default, потом из
; catalog-resolve». Worker не покидает текущую очередь пока в ней есть
; задачи — поэтому mail-sync гарантированно выигрывает голодание.

[program:mzcorp-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mzcorp/artisan queue:work --queue=mail-sync,default,catalog-resolve --sleep=3 --tries=3 --max-time=3600 --max-jobs=200 --memory=600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/mzcorp/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=60
