# deploy/supervisor/mzcorp-worker.conf
; mzcorp воркеры — Laravel queue:work.
;
; Прод: /etc/supervisor/conf.d/mzcorp-worker.conf
; После правки: sudo supervisorctl reread && sudo supervisorctl update
;                && sudo supervisorctl restart mzcorp-worker:* mzcorp-catalog-worker:*
;
; ДВА ИЗОЛИРОВАННЫХ ПУЛА (post-mortem 2026-06-15, см. MEMORY):
;   mzcorp-worker          — mail-sync + default (почтовый конвейер).
;   mzcorp-catalog-worker  — ТОЛЬКО catalog-resolve.
;
; Почему разнесли (а не один пул с приоритетом очередей):
; приоритет --queue решает лишь ЧТО взять свободному воркеру, но НЕ
; вытесняет уже выполняющуюся джобу. ResolvePendingChunkJob под OpenAI
; rate-limit (429) виснет до своего timeout=600с; когда все воркеры
; общего пула разом захватывали такие джобы, mail-sync+default
; замораживались на ~10 мин и почта копилась (события 2026-06-15
; 00:14 / 15:11 / 15:40 / 16:29). Отдельный процесс под catalog-resolve
; гарантирует: каталожные таймауты физически НЕ могут занять почтовые
; воркеры.
;
; Очереди:
;   mail-sync       — IMAP SyncMailboxFolderJob (time-critical).
;   default         — Parse, Deliver, Route, KB resolve, OutgoingQuoteParse,
;                     нотификации, AI-decisions.
;   catalog-resolve — тяжёлая catalog/LLM/pgvector-работа
;                     (ResolvePendingChunkJob и dispatcher).
;
; Порядок --queue для почтового пула: mail-sync раньше default —
; «сначала забери всё из mail-sync». Worker не покидает текущую очередь
; пока в ней есть задачи, поэтому mail-sync не голодает.

[program:mzcorp-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mzcorp/artisan queue:work --queue=mail-sync,default --sleep=3 --tries=3 --max-time=3600 --max-jobs=200 --memory=600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/mzcorp/worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=60

[program:mzcorp-catalog-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/mzcorp/artisan queue:work --queue=catalog-resolve --sleep=5 --tries=1 --max-time=3600 --max-jobs=100 --memory=600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/mzcorp/catalog-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopwaitsecs=620
